# appps/hard/camera.timer
[Unit]
Description=camera timer

[Timer]
OnBootSec=15
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
